# archive/legacy-review/idefix-backoffice/gstech/gstech/packages/brandserver-backend/cron.d/olaspill
0 */5 * * * * /exec cron/update-leaderboard
0 */15 * * * * /exec cron/update-jackpots
