# archive/legacy-review/idefix-backoffice/gstech/gstech/packages/brandserver-backend/cron.d/vie
0 */3 * * * * /exec cron/update-leaderboard
0 */15 * * * * /exec cron/update-jackpots
